# services/linux/collatz.service
[Unit]
Description=Collatz Frontier Coordinator
Documentation=https://github.com/huggablehacker/Collatz-Frontier
After=network.target
Wants=network-online.target

[Service]
# ── Change these to match your setup ─────────────────────────
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/collatz

# Full path to gunicorn (find with: which gunicorn)
# ExecStart for plain HTTP (no TLS):
ExecStart=/home/YOUR_USERNAME/.local/bin/gunicorn \
          -w 1 \
          --threads 32 \
          -b 0.0.0.0:5555 \
          collatz_coordinator:app

# ExecStart for DIRECT TLS (certbot certs, no nginx) — uncomment and edit:
# ExecStart=/home/YOUR_USERNAME/.local/bin/gunicorn \
#           -w 1 \
#           --threads 32 \
#           --certfile /etc/letsencrypt/live/YOUR_DOMAIN/fullchain.pem \
#           --keyfile  /etc/letsencrypt/live/YOUR_DOMAIN/privkey.pem \
#           -b 0.0.0.0:443 \
#           collatz_coordinator:app

# ExecStart for NGINX REVERSE PROXY (certbot via nginx) — uncomment:
# ExecStart=/home/YOUR_USERNAME/.local/bin/gunicorn \
#           -w 1 \
#           --threads 32 \
#           -b 127.0.0.1:5555 \
#           collatz_coordinator:app
# ─────────────────────────────────────────────────────────────

# HTTPS mode — uncomment the relevant line:
# For nginx reverse proxy:
# Environment="COLLATZ_PROXY=1"
# Environment="COLLATZ_DOMAIN=your.domain.com"

# For direct TLS (alternative to gunicorn flags above):
# Environment="COLLATZ_CERTFILE=/etc/letsencrypt/live/YOUR_DOMAIN/fullchain.pem"
# Environment="COLLATZ_KEYFILE=/etc/letsencrypt/live/YOUR_DOMAIN/privkey.pem"

# GitHub token for nightly uploads:
# Environment="GITHUB_TOKEN=ghp_your_token_here"

Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
KillMode=mixed
StandardOutput=journal
StandardError=journal
SyslogIdentifier=collatz

[Install]
WantedBy=multi-user.target
